# systemd/alarm.timer
[Unit]
Description=Daily

[Timer]
OnCalendar=*-*-* 13:00:00
Persistent=false
 
[Install]
WantedBy=timers.target
